# app/crontab
# minute hour day(of month) month(day of week) command
# This is take a photo on the hour for hours 11, 15, 19, 23 UTC
0 11,15,19,23 * * * python /app/snap.py > /proc/1/fd/1 2>/proc/1/fd/2